# zangetsu/deploy/systemd/arena23-orchestrator.service
[Unit]
Description=Zangetsu V5 Arena23 Orchestrator
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=j13
WorkingDirectory=/home/j13/j13-ops
EnvironmentFile=/etc/zangetsu/zangetsu.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/j13/j13-ops/zangetsu/.venv/bin/python3 /home/j13/j13-ops/zangetsu/services/arena23_orchestrator.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
