# zangetsu/deploy/systemd/console-api.service
[Unit]
Description=Zangetsu V5 Console API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=j13
WorkingDirectory=/home/j13/j13-ops
EnvironmentFile=/home/j13/j13-ops/zangetsu/secret/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/j13/j13-ops/zangetsu/.venv/bin/uvicorn zangetsu.console.run:app --app-dir /home/j13/j13-ops --host 0.0.0.0 --port 9900
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
