# zangetsu/deploy/systemd/live-trader.service
[Unit]
Description=Zangetsu V5 Live Trader
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=j13
WorkingDirectory=/home/j13/j13-ops
EnvironmentFile=/home/j13/j13-ops/zangetsu/secret/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/j13/j13-ops/zangetsu/.venv/bin/python3 /home/j13/j13-ops/zangetsu/services/live_trader.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
